# YourBot/deployment/discord-health.service
[Unit]
Description=Discord Bot Health Check Server
After=network.target

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/opt/discord-bot
ExecStart=/bin/bash -c 'mkdir -p /var/log/discord-bot && /usr/bin/node health-check.js >> /var/log/discord-bot/health.log 2>&1'
Restart=always
RestartSec=5
Environment=NODE_ENV=production
SyslogIdentifier=discord-health

[Install]
WantedBy=multi-user.target