# YourBot/deployment/discord-health.service
[Unit]
Description=Discord Bot Health Check Server
After=network.target

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/opt/discord-bot
ExecStart=/usr/bin/node health-check.js
Restart=always
RestartSec=5
Environment=NODE_ENV=production
StandardOutput=journal
StandardError=journal
SyslogIdentifier=discord-health

[Install]
WantedBy=multi-user.target